# ansible/roles/ncmpcpp/files/mopidy.service
[Unit]
Description=mopidy
After=network.target

[Service]
Type=oneshot
ExecStart=/home/johan/.local/bin/mopidy

[Install]
WantedBy=default.target
